# backend/deploy/aqua-backend.service
[Unit]
Description=Aqua Monitoring FastAPI Backend
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/aqua/backend
EnvironmentFile=/var/www/aqua/backend/.env
ExecStart=/var/www/aqua/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
